# files/etc/s6-overlay/s6-rc.d/migrate/up
#!/command/execlineb -P
/bin/sh /usr/local/bin/migrate.sh
